# telegram_bot_service.service
[Unit]
Description=REO Dashboard Telegram Bot
After=network.target

[Service]
Type=simple
User=graph
WorkingDirectory=/home/graph/ftpbox/reo
ExecStart=/usr/bin/python3 /home/graph/ftpbox/reo/telegram_bot.py
Restart=always
RestartSec=10
StandardOutput=append:/home/graph/ftpbox/reo/telegram_bot.log
StandardError=append:/home/graph/ftpbox/reo/telegram_bot.log

# Environment file (optional, if you want to load .env through systemd)
# EnvironmentFile=/home/graph/ftpbox/reo/.env

[Install]
WantedBy=multi-user.target

